# contrib/it-asset-manager.service
[Unit]
Description=IT Asset Manager (Flask/Gunicorn)
After=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
Environment="PATH=__APP_DIR__/.venv/bin"
Environment="FLASK_DEBUG=0"
Environment="SECRET_KEY=__CHANGE_ME__"
ExecStart=__APP_DIR__/.venv/bin/gunicorn -c gunicorn.conf.py "server:app"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.targ